# deploy/systemd/group-member-guardian.service
[Unit]
Description=Group Member Guardian Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Change these to your deployment user/group
User=root
Group=root

WorkingDirectory=/root/vondo-spirit-tgbot
EnvironmentFile=/root/vondo-spirit-tgbot/.env
ExecStart=/root/vondo-spirit-tgbot/.venv/bin/python -m bot.main

Restart=always
RestartSec=5
TimeoutStopSec=20

# Basic hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full

[Install]
WantedBy=multi-user.target
